# ntt.service
[Unit]
Description=NTTBot
After=syslog.target network.target

[Service]
Type=simple

WorkingDirectory=/usr/local/ntt/bin
ExecStart=/bin/bash /usr/local/ntt/bin/start.sh

User=root
Group=root

[Install]

WantedBy=multi-user.target
