# ntt.service
[Unit]
Description=sell
After=syslog.target network.target

[Service]
Type=forking

WorkingDirectory=/bot
ExecStart=/bin/base /bot/start.sh

User=root
Group=root

[Install]

WantedBy=multi-user.target